# etc/supervisor/blitzka-prod.conf
[unix_http_server]
file=/var/run/supervisor.sock
chmod=0700

[supervisord]
logfile=/opt/blitzka/data/logs/supervisord.log
environment=LC_ALL='en_US.UTF-8',LANG='en_US.UTF-8',PYTHONIOENCODING='UTF-8'

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///var/run/supervisor.sock

[program:nginx]
priority=2
command=/usr/sbin/nginx -c /opt/blitzka/etc/nginx/blitzka-prod.conf
autorestart=true
startretries=3
stopwaitsecs=10
stopasgroup=true
stderr_logfile=/opt/blitzka/data/logs/nginx.err.log
stdout_logfile=/opt/blitzka/data/logs/nginx.out.log
user=root

[program:postgres]
priority=3
command=/usr/lib/postgresql/10/bin/postgres -D /opt/blitzka/data/database
autorestart=true
startretries=3
stopwaitsecs=10
stopasgroup=true
stderr_logfile=/opt/blitzka/data/logs/postgres.log
stdout_logfile=/opt/blitzka/data/logs/postgres.log
user=postgres

[program:redis]
priority=3
command=/usr/bin/redis-server /etc/redis/redis.conf --unixsocket /tmp/redis.sock --daemonize no
autorestart=true
startretries=3
stopwaitsecs=30
stopasgroup=true
stderr_logfile=/opt/blitzka/data/logs/redis.log
stdout_logfile=/opt/blitzka/data/logs/redis.log
user=redis

[program:daphne]
priority=4
command=/opt/blitzka/grater-django/venv/bin/daphne -b 127.0.0.1 -p 8000 grater.asgi:channel_layer
directory=/opt/blitzka/grater-django
autorestart=true
startretries=3
stopwaitsecs=3
stopasgroup=true
stderr_logfile=/opt/blitzka/data/logs/daphne.log
stdout_logfile=/opt/blitzka/data/logs/daphne.log
environment=DJANGO_SETTINGS_MODULE="grater.settings",ODDSLINGERS_ENV='PROD',PATH="/opt/blitzka/grater-django/venv/bin:%(ENV_PATH)s",LANG="en_US.UTF-8",LC_ALL="en_US.UTF-8"
user=www-data

[program:dramatiq]
priority=5
command=nice -15 /opt/blitzka/grater-django/venv/bin/python manage.py rundramatiq --processes 1 --threads 2 --no-reload
directory=/opt/blitzka/grater-django
autorestart=true
startretries=3
stopwaitsecs=3
stopasgroup=true
stderr_logfile=/opt/blitzka/data/logs/dramatiq.log
stdout_logfile=/opt/blitzka/data/logs/dramatiq.log
environment=DJANGO_SETTINGS_MODULE="grater.settings",ODDSLINGERS_ENV='PROD',PATH="/opt/blitzka/grater-django/venv/bin:%(ENV_PATH)s",LANG="en_US.UTF-8",LC_ALL="en_US.UTF-8",PYTHONIOENCODING='UTF-8'
user=www-data

[program:yacron]
priority=10
command=/opt/blitzka/grater-django/venv/bin/yacron -c /opt/blitzka/etc/yacron/blitzka-prod.yaml
autostart=true
autorestart=true
autorestart=true
startretries=3
stopwaitsecs=3
stopasgroup=true
stderr_logfile=/opt/blitzka/data/logs/yacron.log
stdout_logfile=/opt/blitzka/data/logs/yacron.log
environment=DJANGO_SETTINGS_MODULE="grater.settings",ODDSLINGERS_ENV='PROD',PATH="/opt/blitzka/grater-django/venv/bin:%(ENV_PATH)s",LANG="en_US.UTF-8",LC_ALL="en_US.UTF-8"
user=root

[program:django-http-worker]
priority=6
process_name=django-http-worker-%(process_num)s
numprocs=2
command=/opt/blitzka/grater-django/venv/bin/python manage.py runworker --exclude-channels=websocket.*
directory=/opt/blitzka/grater-django
autorestart=true
startretries=3
stopasgroup=true
stderr_logfile=/opt/blitzka/data/logs/http-worker.log
stdout_logfile=/opt/blitzka/data/logs/http-worker.log
environment=DJANGO_SETTINGS_MODULE="grater.settings",ODDSLINGERS_ENV='PROD',PATH="/opt/blitzka/grater-django/venv/bin:%(ENV_PATH)s",LANG="en_US.UTF-8",LC_ALL="en_US.UTF-8"
user=www-data

[program:django-socket-worker]
priority=7
process_name=django-socket-worker-%(process_num)s
numprocs=4
command=/opt/blitzka/grater-django/venv/bin/python manage.py runworker --only-channels=websocket.*
directory=/opt/blitzka/grater-django
autorestart=true
startretries=3
stopasgroup=true
stderr_logfile=/opt/blitzka/data/logs/socket-worker.log
stdout_logfile=/opt/blitzka/data/logs/socket-worker.log
environment=DJANGO_SETTINGS_MODULE="grater.settings",ODDSLINGERS_ENV='PROD',PATH="/opt/blitzka/grater-django/venv/bin:%(ENV_PATH)s",LANG="en_US.UTF-8",LC_ALL="en_US.UTF-8"
user=www-data

[program:botbeat]
priority=8
command=nice -10 /opt/blitzka/grater-django/venv/bin/python manage.py bot_heartbeat
directory=/opt/blitzka/grater-django
autorestart=true
startretries=3
stopwaitsecs=3
stopasgroup=true
stderr_logfile=/opt/blitzka/data/logs/botbeat.log
stdout_logfile=/opt/blitzka/data/logs/botbeat.log
environment=DJANGO_SETTINGS_MODULE="grater.settings",ODDSLINGERS_ENV='PROD',PATH="/opt/blitzka/grater-django/venv/bin:%(ENV_PATH)s",LANG="en_US.UTF-8",LC_ALL="en_US.UTF-8"
user=www-data

[program:support-bot]
priority=9
command=nice -10 /opt/blitzka/grater-django/venv/bin/python manage.py support_bot
directory=/opt/blitzka/grater-django
autorestart=true
startretries=3
stopwaitsecs=3
stopasgroup=true
stderr_logfile=/opt/blitzka/data/logs/support_bot.log
stdout_logfile=/opt/blitzka/data/logs/support_bot.log
environment=DJANGO_SETTINGS_MODULE="grater.settings",ODDSLINGERS_ENV='PROD',PATH="/opt/blitzka/grater-django/venv/bin:%(ENV_PATH)s",LANG="en_US.UTF-8",LC_ALL="en_US.UTF-8"
user=www-data


[group:blitzka-base]
programs=nginx,postgres,redis

[group:blitzka-bg]
programs=yacron,dramatiq,support-bot,botbeat

[group:blitzka-django]
programs=daphne,django-http-worker,django-socket-worker
