# etc/supervisor/oddslingers-beta.conf
[unix_http_server]
file=/var/run/supervisor.sock
chmod=0700

[supervisord]
logfile=/opt/oddslingers.poker/data/logs/supervisord.log

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///var/run/supervisor.sock

[program:sshfs]
priority=1
command=/usr/bin/sshfs -o allow_other -o reconnect -f lime:/tank/oddslingers-beta/data/support_tickets /data/support_tickets
autorestart=true
startretries=3
stopwaitsecs=300
stopasgroup=true
stderr_logfile=/opt/oddslingers.poker/data/logs/sshfs.log
stdout_logfile=/opt/oddslingers.poker/data/logs/sshfs.log
user=root

[program:nginx]
priority=2
command=/usr/sbin/nginx -c /opt/oddslingers.poker/etc/nginx/oddslingers-beta.conf
autorestart=true
startretries=3
stopwaitsecs=10
stopasgroup=true
stderr_logfile=/opt/oddslingers.poker/data/logs/nginx.err.log
stdout_logfile=/opt/oddslingers.poker/data/logs/nginx.out.log
user=root

[program:postgres]
priority=3
command=/usr/lib/postgresql/10/bin/postgres -D /opt/oddslingers.poker/data/database
autorestart=true
startretries=3
stopwaitsecs=10
stopasgroup=true
stderr_logfile=/opt/oddslingers.poker/data/logs/postgres.log
stdout_logfile=/opt/oddslingers.poker/data/logs/postgres.log
user=postgres

[program:redis]
priority=3
command=/usr/bin/redis-server /etc/redis/redis.conf --unixsocket /tmp/redis.sock --daemonize no
autorestart=true
startretries=3
stopwaitsecs=30
stopasgroup=true
stderr_logfile=/opt/oddslingers.poker/data/logs/redis.log
stdout_logfile=/opt/oddslingers.poker/data/logs/redis.log
user=redis

[program:daphne]
priority=4
command=/opt/oddslingers.poker/core/venv/bin/daphne -b 127.0.0.1 -p 8000 oddslingers.asgi:channel_layer
directory=/opt/oddslingers.poker/core
autorestart=true
startretries=3
stopwaitsecs=3
stopasgroup=true
stderr_logfile=/opt/oddslingers.poker/data/logs/daphne.log
stdout_logfile=/opt/oddslingers.poker/data/logs/daphne.log
environment=DJANGO_SETTINGS_MODULE="oddslingers.settings",ODDSLINGERS_ENV='BETA',PATH="/opt/oddslingers.poker/core/venv/bin:%(ENV_PATH)s",LANG="en_US.UTF-8",LC_ALL="en_US.UTF-8"
user=www-data

[program:dramatiq]
priority=5
command=nice -15 /opt/oddslingers.poker/core/venv/bin/python manage.py rundramatiq --processes 1 --threads 2 --no-reload
directory=/opt/oddslingers.poker/core
autorestart=true
startretries=3
stopwaitsecs=3
stopasgroup=true
stderr_logfile=/opt/oddslingers.poker/data/logs/dramatiq.log
stdout_logfile=/opt/oddslingers.poker/data/logs/dramatiq.log
environment=DJANGO_SETTINGS_MODULE="oddslingers.settings",ODDSLINGERS_ENV='BETA',PATH="/opt/oddslingers.poker/core/venv/bin:%(ENV_PATH)s",LANG="en_US.UTF-8",LC_ALL="en_US.UTF-8"
user=www-data

[program:yacron]
priority=10
command=/opt/oddslingers.poker/core/venv/bin/yacron -c /opt/oddslingers.poker/etc/yacron/oddslingers-beta.yaml
autostart=true
autorestart=true
autorestart=true
startretries=3
stopwaitsecs=3
stopasgroup=true
stderr_logfile=/opt/oddslingers.poker/data/logs/yacron.log
stdout_logfile=/opt/oddslingers.poker/data/logs/yacron.log
environment=DJANGO_SETTINGS_MODULE="oddslingers.settings",ODDSLINGERS_ENV='BETA',PATH="/opt/oddslingers.poker/core/venv/bin:%(ENV_PATH)s",LANG="en_US.UTF-8",LC_ALL="en_US.UTF-8"
user=root

[program:django-http-worker]
priority=6
process_name=django-http-worker-%(process_num)s
numprocs=2
command=/opt/oddslingers.poker/core/venv/bin/python manage.py runworker --exclude-channels=websocket.*
directory=/opt/oddslingers.poker/core
autorestart=true
startretries=3
stopasgroup=true
stderr_logfile=/opt/oddslingers.poker/data/logs/http-worker.log
stdout_logfile=/opt/oddslingers.poker/data/logs/http-worker.log
environment=DJANGO_SETTINGS_MODULE="oddslingers.settings",ODDSLINGERS_ENV='BETA',PATH="/opt/oddslingers.poker/core/venv/bin:%(ENV_PATH)s",LANG="en_US.UTF-8",LC_ALL="en_US.UTF-8"
user=www-data

[program:django-socket-worker]
priority=7
process_name=django-socket-worker-%(process_num)s
numprocs=4
command=/opt/oddslingers.poker/core/venv/bin/python manage.py runworker --only-channels=websocket.*
directory=/opt/oddslingers.poker/core
autorestart=true
startretries=3
stopasgroup=true
stderr_logfile=/opt/oddslingers.poker/data/logs/socket-worker.log
stdout_logfile=/opt/oddslingers.poker/data/logs/socket-worker.log
environment=DJANGO_SETTINGS_MODULE="oddslingers.settings",ODDSLINGERS_ENV='BETA',PATH="/opt/oddslingers.poker/core/venv/bin:%(ENV_PATH)s",LANG="en_US.UTF-8",LC_ALL="en_US.UTF-8"
user=www-data

[program:botbeat]
priority=8
command=nice -10 /opt/oddslingers.poker/core/venv/bin/python manage.py bot_heartbeat
directory=/opt/oddslingers.poker/core
autorestart=true
startretries=3
stopwaitsecs=3
stopasgroup=true
stderr_logfile=/opt/oddslingers.poker/data/logs/botbeat.log
stdout_logfile=/opt/oddslingers.poker/data/logs/botbeat.log
environment=DJANGO_SETTINGS_MODULE="oddslingers.settings",ODDSLINGERS_ENV='BETA',PATH="/opt/oddslingers.poker/core/venv/bin:%(ENV_PATH)s",LANG="en_US.UTF-8",LC_ALL="en_US.UTF-8"
user=www-data

# [program:support-bot]
# priority=9
# command=nice -10 /opt/oddslingers.poker/core/venv/bin/python manage.py support_bot
# directory=/opt/oddslingers.poker/core
# autorestart=true
# startretries=3
# stopwaitsecs=3
# stopasgroup=true
# stderr_logfile=/opt/oddslingers.poker/data/logs/support_bot.log
# stdout_logfile=/opt/oddslingers.poker/data/logs/support_bot.log
# environment=DJANGO_SETTINGS_MODULE="oddslingers.settings",ODDSLINGERS_ENV='BETA',PATH="/opt/oddslingers.poker/core/venv/bin:%(ENV_PATH)s",LANG="en_US.UTF-8",LC_ALL="en_US.UTF-8"
# user=www-data

# [program:postgres-tunnel]
# priority=8
# command=/usr/bin/stunnel /opt/oddslingers.poker/etc/stunnel/postgres-server.conf
# autorestart=true
# startretries=3
# stopwaitsecs=10
# stopasgroup=true
# stderr_logfile=/opt/oddslingers.poker/data/logs/postgres-tunnel.log
# stdout_logfile=/opt/oddslingers.poker/data/logs/postgres-tunnel.log
# user=root

# [program:redis-tunnel]
# priority=8
# command=/usr/bin/stunnel /opt/oddslingers.poker/etc/stunnel/redis-server.conf
# autorestart=true
# startretries=3
# stopwaitsecs=30
# stopasgroup=true
# stderr_logfile=/opt/oddslingers.poker/data/logs/redis-tunnel.log
# stdout_logfile=/opt/oddslingers.poker/data/logs/redis-tunnel.log
# user=root

[group:oddslingers-base]
programs=sshfs,nginx,postgres,redis

[group:oddslingers-bg]
programs=yacron,dramatiq,botbeat

[group:oddslingers-django]
programs=daphne,django-http-worker,django-socket-worker
