# etc/supervisor/ubuntu-dev.conf
[unix_http_server]
file=/var/run/supervisor.sock
chmod=0700

[supervisord]
logfile=/opt/oddslingers.poker/data/logs/supervisord.log
environment=LC_ALL='en_US.UTF-8',LANG='en_US.UTF-8'

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///var/run/supervisor.sock

[program:nginx]
priority=1
command=/usr/sbin/nginx -c /opt/oddslingers.poker/etc/nginx/oddslingers-dev.conf
autorestart=true
startretries=3
stopwaitsecs=10
stopasgroup=true
stderr_logfile=/opt/oddslingers.poker/data/logs/nginx.err.log
stdout_logfile=/opt/oddslingers.poker/data/logs/nginx.out.log

[program:postgres]
user=postgres
priority=2
command=/usr/lib/postgresql/12/bin/postgres -D /opt/oddslingers.poker/data/postgres
autorestart=true
startretries=3
stopwaitsecs=10
stopasgroup=true
stderr_logfile=/opt/oddslingers.poker/data/logs/postgres.err.log
stdout_logfile=/opt/oddslingers.poker/data/logs/postgres.out.log

[program:redis]
priority=2
command=/usr/bin/redis-server /etc/redis/redis.conf --unixsocket /tmp/redis.sock --daemonize no
autorestart=true
startretries=3
stopwaitsecs=30
stopasgroup=true
stderr_logfile=/opt/oddslingers.poker/data/logs/redis.err.log
stdout_logfile=/opt/oddslingers.poker/data/logs/redis.out.log

[program:daphne]
priority=4
command=/opt/oddslingers.poker/core/.venv/bin/daphne -b 127.0.0.1 -p 8000 oddslingers.asgi:channel_layer
directory=/opt/oddslingers.poker/core
autostart=false
autorestart=true
startretries=3
stopwaitsecs=3
stopasgroup=true
stderr_logfile=/opt/oddslingers.poker/data/logs/daphne.err.log
stdout_logfile=/opt/oddslingers.poker/data/logs/daphne.out.log
environment=DJANGO_SETTINGS_MODULE="oddslingers.settings",ODDSLINGERS_ENV='DEV',PATH="/opt/oddslingers.poker/core/.venv/bin:%(ENV_PATH)s"

[program:dramatiq]
priority=5
command=nice -15 /opt/oddslingers.poker/core/.venv/bin/python manage.py rundramatiq --processes 1 --threads 1
directory=/opt/oddslingers.poker/core
autostart=false
autorestart=true
startretries=3
stopwaitsecs=3
stopasgroup=true
stderr_logfile=/opt/oddslingers.poker/data/logs/dramatiq.err.log
stdout_logfile=/opt/oddslingers.poker/data/logs/dramatiq.out.log
environment=DJANGO_SETTINGS_MODULE="oddslingers.settings",ODDSLINGERS_ENV='DEV',PATH="/opt/oddslingers.poker/core/.venv/bin:%(ENV_PATH)s"

[program:django]
priority=6
command=/opt/oddslingers.poker/core/.venv/bin/python manage.py runworker
directory=/opt/oddslingers.poker/core
autostart=false
autorestart=true
startretries=3
stopasgroup=true
stderr_logfile=/opt/oddslingers.poker/data/logs/http-worker.err.log
stdout_logfile=/opt/oddslingers.poker/data/logs/http-worker.out.log
environment=DJANGO_SETTINGS_MODULE="oddslingers.settings",ODDSLINGERS_ENV='DEV',PATH="/opt/oddslingers.poker/core/.venv/bin:%(ENV_PATH)s",

[program:botbeat]
priority=10
command=nice -10 /opt/oddslingers.poker/core/.venv/bin/python manage.py bot_heartbeat --stupid
directory=/opt/oddslingers.poker/core
autostart=false
autorestart=true
startretries=3
stopwaitsecs=3
stopasgroup=true
stderr_logfile=/opt/oddslingers.poker/data/logs/botbeat.err.log
stdout_logfile=/opt/oddslingers.poker/data/logs/botbeat.out.log
environment=DJANGO_SETTINGS_MODULE="oddslingers.settings",ODDSLINGERS_ENV='DEV',PATH="/opt/oddslingers.poker/core/.venv/bin:%(ENV_PATH)s"

[program:yacron]
priority=10
command=/opt/oddslingers.poker/core/.venv/bin/yacron -c /opt/oddslingers.poker/etc/yacron/oddslingers-dev.yaml
autostart=true
autorestart=true
autorestart=true
startretries=3
stopwaitsecs=3
stopasgroup=true
stderr_logfile=/opt/oddslingers.poker/data/logs/yacron.log
stdout_logfile=/opt/oddslingers.poker/data/logs/yacron.log
environment=DJANGO_SETTINGS_MODULE="oddslingers.settings",ODDSLINGERS_ENV='DEV',PATH="/opt/oddslingers.poker/core/.venv/bin:%(ENV_PATH)s",LANG="en_US.UTF-8",LC_ALL="en_US.UTF-8"
user=root


[group:oddslingers.poker]
programs=nginx,postgres,redis

[group:oddslingers-bg]
programs=yacron,dramatiq,botbeat

[group:oddslingers-django]
programs=daphne,django
